# conf/logstash.service
[Unit]
Description=Logstash
Documentation=https://www.elastic.co/products/logstash
After=network.target
#ConditionPathExists=/etc/logstash.conf

[Service]
User=deploy
Group=deploy
ExecStart=/home/ec2-user/ES-Tutorial-ELK/packages/logstash/bin/logstash

[Install]
WantedBy=multi-user.target
